# ansible/newsbot.service
[Service]
Restart=on-failure
WorkingDirectory=/home/ec2-user/newsbot
ExecStart=/usr/bin/python3 run start
EnvironmentFile=/home/ec2-user/bot/.env
