# ansible/newsbot.service
[Service]
Restart=on-failure
WorkingDirectory=/home/ec2-user/bot
ExecStart=/usr/bin/python3 gear_news_bot
EnvironmentFile=/home/ec2-user/bot/.env
